# project/nodes/attacker2/systemd-service/attacker2.path
[Path]
PathExists=/opt/attacker2/control/restart_required

[Install]
WantedBy=multi-user.target